# supervisor-tuya-listener.conf
[program:tuya-listener]
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/mnt/data/code/mcp
startretries=999                ; max # of serial start failures when starting (default 3)
startsecs=3                   ; # of secs prog must stay up to be running (def. 1)
autorestart=true        ; when to restart if exited after running (def: unexpected)
redirect_stderr=true
stdout_logfile=./log/%(program_name)s.log
stdout_logfile_maxbytes=50MB
command=/bin/bash -c "/mnt/data/code/mcp/.venv/bin/python -m mcp.controllers.tuya_listener"


; command=/bin/bash -c "/usr/local/bin/docker-compose up --force-recreate --build tuya-listener"
        ; "/usr/local/bin/docker-compose logs --timestamps --follow tuya-listener"
        ; "tail -n40 -f ./log/%(program_name)s.log"

